# lesson2.bash.scripts/3/lesson2.service
[Unit]
Description=lesson2 empty echo

[Service]
ExecStart=/root/devops/lesson2/3/echo.sh 

[Install]
WantedBy=multi-user.target 
